# etc/supervisor/conf.d/gesangprasojo.ini
[program:brain_gesang_prasojo]
command=php -S localhost:8000 -t /home/app/www_brain/www_backend/public
autostart=true
autorestart=true
user=app